# private_dot_config/systemd/user/low-battery-warn.service
[Unit]
Description=Low Battery Warning

[Service]
Type=oneshot
ExecStart=%h/.local/bin/low-battery-warn

