# private_dot_config/systemd/user/mako-dnd-off.service
[Unit]
Description=Enable Mako do-not-disturb

[Service]
Type=oneshot
ExecStart=/usr/bin/makoctl mode -r do-not-disturb
